# openwrt-2.4.2011-trunk/target/linux/g2/base-files/etc/init.d/zzcomplete
#!/bin/sh /etc/rc.common

START=99

start() {
  [ -x /bin/psplash-write ] && /bin/psplash-write "MSG Ready!"
}

