# openwrt-2.4.2011-trunk/target/linux/g2/base-files/etc/inittab
# 
# /etc/inittab
#
# $Id: inittab,v 1.1.1.1 2011/02/07 19:20:42 peebles Exp $
#

::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K stop

::respawn:/sbin/getty -L ttyS0 115200 vt100

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcS K stop

